# install-scripts/systemd-services/kanata-hidden.service
[Unit]
Description=Kanata Service

[Service]
ExecStart=/usr/bin/kanata -c %h/.config/kanata/kanata.kbd
Restart=no

[Install]
WantedBy=default.target

